# templates/default/upstart.conf.erb
description "upstart configuration for openwatch django gunicorn"

setuid <%= @service_user %>

start on net-device-up
stop on shutdown

respawn

script
  . <%= @virtualenv_path %>/bin/activate
  cd <%= @app_root %>/<%= @app_name %>
  exec gunicorn <%= @app_name %>.wsgi:application -w 8 -u <%= @service_user %> --log-level=debug --access-logfile=<%= @access_log_path %> --error-logfile=<%= @error_log_path %>
end script