# templates/default/upstart.conf.erb
#
# /etc/init/ow_media_capture.conf
#
# Generated by Chef for <%= node[:fqdn] %>
#

description "OpenWatch NodeMediaCapture application"

setuid <%= @service_user %>

start on net-device-up
stop on shutdown

respawn

script
  cd <%= @app_root %>
  exec <%= @app_root %>/<%= @run_script %> >> <%= @log_path %> 2>&1
end script
